# systemd/discord-bot.service
# File 1: discord-bot.service
# Location: systemd/discord-bot.service
# Purpose: Tells systemd how to run your bot
[Unit]
Description=Discord Bot Service
After=network.target

[Service]
Type=simple
User=discordbot
EnvironmentFile=/etc/discord-bot.env
ExecStart=/usr/local/bin/discord-bot-runner.sh
WorkingDirectory=/opt/discord-bot/embed-fixer
Restart=always
RestartSec=3
StandardOutput=append:/var/log/discord-bot.log
StandardError=append:/var/log/discord-bot.error.log

# Security settings
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ReadWritePaths=/opt/discord-bot

[Install]
WantedBy=multi-user.target